# hass
#!/sbin/openrc-run

# OpenRC service file for Home Assistant in a Python virtual environment
#
# Config variables
# H_USER: Run as user
# H_GROUP: Run as group
# PYTHON_VENV_PATH: Path to python virtual environment

command="hass"
pidfile="/run/${RC_SVCNAME}/${RC_SVCNAME}.pid"
command_args="--pid-file ${pidfile}"
command_args_background="--daemon"
command_user="${H_USER}:${H_GROUP}"

extra_commands="update"

# Home Assistant takes a while to terminate. Wait 20s
retry="20"

quit() {
  eerror "$@"
  exit 1
}

start_pre() {
  [ -z "${H_USER}" ] && quit "Must specify user!"
  [ -z "${H_GROUP}" ] && quit "Must specify group!"

  [ "${H_USER}" = "root" ] &&  quit "User root not allowed!"
  [ "${H_GROUP}" = "root" ] && quit "Group root not allowed!"

  [ -f ${PYTHON_VENV_PATH}/bin/activate ] || \
    quit "Cannot source venv ${PYTHON_VENV_PATH}/bin/activate"

  checkpath -q -d -m 0750 -o "${command_user}" "$(dirname "${pidfile}")"
  source ${PYTHON_VENV_PATH}/bin/activate
}

update() {
  [ -z "${H_USER}" ] && quit "Must specify user!"
  [ "${H_USER}" = "root" ] &&  quit "User root not allowed!"

  [ -f ${PYTHON_VENV_PATH}/bin/activate ] || \
    quit "Cannot source venv ${PYTHON_VENV_PATH}/bin/activate"

  if service_stopping || service_starting || service_inactive || \
    service_crashed || service_started; then
    quit "Service must be stopped to update"
  fi

  su -s /bin/bash -c \
    "source ${PYTHON_VENV_PATH}/bin/activate; \
      python3 -m pip install --upgrade homeassistant" \
    ${H_USER}
}
